# src/gitblob.service
[Service]
Environment=PORT=8976
Environment=GIT_ROOT=/var/lib/git
ExecStart=/usr/local/bin/gitblob

[Install]
WantedBy=multi-user.target

